# deploy/service/golden-hour.service
[Unit]
Description=Golden Hour Study Bot (OpenClaw Gateway)
Documentation=https://github.com/margoshkagt-star/Golden-Hour
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=golden-hour
Group=golden-hour
WorkingDirectory=/opt/golden-hour

# --- Environment ---
# Copy deploy/.env.example to /opt/golden-hour/.env (chmod 600, owned by golden-hour)
EnvironmentFile=-/opt/golden-hour/.env

# Non-standard ports — override in /opt/golden-hour/.env
# APP_PORT     = internal health/API port (default 47832)
# WEBHOOK_PORT = Telegram webhook if not using long-poll (must be 443/80/88/8443 or via reverse proxy)
Environment=NODE_ENV=production
Environment=APP_PORT=47832
Environment=WEBHOOK_PORT=47843
Environment=TZ=Europe/Moscow

# --- Run ---
ExecStart=/usr/local/bin/openclaw gateway

# Do NOT use ExecReload — OpenClaw's SIGHUP behaviour is unverified.
# Use 'systemctl restart golden-hour' for all updates.
# ExecReload=/bin/kill -HUP $MAINPID

# --- Restart policy ---
Restart=always
RestartSec=10s
# After 5 crashes within 120s, systemd stops restarting (prevents infinite loop).
# Monitor with: journalctl -u golden-hour -f
# Alert on failure via: systemctl status golden-hour
StartLimitIntervalSec=120s
StartLimitBurst=5

# Graceful shutdown: give in-flight requests time to finish (writes must be
# atomic: write to .tmp then rename() to prevent partial-write corruption).
TimeoutStopSec=30s

# --- Resource limits ---
# Node.js uses several threads (libuv threadpool + GC + main). 64 is too low.
# Profile actual usage; 256 is a safe default for a multi-user Node bot.
MemoryMax=512M
MemoryHigh=400M
TasksMax=256
LimitNOFILE=65536

# --- Hardening ---
NoNewPrivileges=yes
ProtectSystem=strict
# Allow writing only to project dir; PrivateTmp gives its own /tmp already.
# Do NOT add /tmp here — PrivateTmp already provides a private writable /tmp,
# and adding /tmp to ReadWritePaths would re-open the real shared /tmp.
ReadWritePaths=/opt/golden-hour
PrivateTmp=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes

[Install]
WantedBy=multi-user.target
